# jasm.service
[Unit]
Description=JASM Server Daemon

[Service]
Type=forking
ExecStart=/usr/local/bin/jasm

[Install]
WantedBy=multi-user.target